# root/lib/udev/rules.d/90-finnix.rules
ACTION=="add|remove|change", SUBSYSTEM=="block", RUN+="/bin/sh -c '/usr/sbin/rebuildfstab &'"
