# root/lib/finnix/service.d/getty-skel/finish
#!/bin/sh

[ -f ./config ] && . ./config
[ -n "$TTY" ] || exit 1

if [ -x /usr/sbin/finnix-utmp ]; then
  /usr/sbin/finnix-utmp del $TTY
fi
