# root/lib/finnix/service.d/getty-skel/finish
#!/bin/sh

[ -f ./config ] && . ./config
[ -n "$TTY" ] || exit 1

if [ -x /usr/sbin/utmptty ]; then
  /usr/sbin/utmptty del </dev/$TTY 2>/dev/null
fi
