# root/lib/finnix/service.d/getty-skel/run
#!/bin/sh

[ -f ./config ] && . ./config
[ -n "$TTY" ] || exit 1

if [ -x /usr/sbin/finnix-utmp ]; then
  /usr/sbin/finnix-utmp add $TTY root
fi
exec chpst -P /sbin/shell-login $TTY
